# supervisord.ini
[supervisord]
nodaemon=true
loglevel=debug

[program:dnscrypt-proxy]
command=dnscrypt-proxy
autostart=true
autorestart=true
user=root

[program:cron]
command = cron -f -L 15
autostart=true
autorestart=true
user=root
