# supervisord.ini
[supervisord]
nodaemon=true
loglevel=debug
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:dnscrypt-proxy]
command=dnscrypt-proxy --config=/etc/dnscrypt-proxy.toml
autostart=true
autorestart=true
user=root
stdout_events_enabled = true
stderr_events_enabled = true

[program:cron]
command = cron -f -L 15
autostart=true
autorestart=true
user=root
stdout_events_enabled = true
stderr_events_enabled = true

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler
